# shell/history/1.crontab.example
5 0 * * * /home/ec2-user/history/currency.sh >> /home/ec2-user/history/currency.log 2>&1
