# shell/history/2.crontab.example
*/5 * * * * /home/trade/run/cleanlog.sh 10 >> /home/trade/a
